# systemd/sufe-8021x.service
[Unit]
Description=SUFE campus 802.1X authentication
After=network-pre.target
Before=network-online.target
Wants=network-pre.target

[Service]
Type=simple
EnvironmentFile=/etc/sufe-8021x/sufe-8021x.env
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/usr/local/sbin/sufe-campus-net.sh
ExecStartPre=/bin/rm -f /run/sufe-8021x-authenticated
ExecStart=/usr/bin/python3 /usr/local/sbin/sufe_8021x_auth.py --env-file /etc/sufe-8021x/sufe-8021x.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
